# examples/cron.d/cacti
# Example of how to run this with your poller.php script
*/1 * * * * www-data php /usr/share/cacti/poller.php > /dev/null 2>&1 && /path/to/weathermap-export.sh

# Generate time-lapse automatically on first day of the month (00:00)
#0 0 1 * * www-data /path/to/timelapse-export.py -f -m "/path/to/stored/frames" "/path/to/timelapse.mkv"
